# scripts/com.tiktok-affiliate.post.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tiktok-affiliate.post</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/daichi/Desktop/dev/claude-code/tiktok-shop-affiliate/venv/bin/python</string>
        <string>/Users/daichi/Desktop/dev/claude-code/tiktok-shop-affiliate/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/daichi/Desktop/dev/claude-code/tiktok-shop-affiliate</string>

    <key>StartCalendarInterval</key>
    <array>
        <!-- 8:00 -->
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 14:00 -->
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 20:00 -->
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/daichi/Desktop/dev/claude-code/tiktok-shop-affiliate/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/daichi/Desktop/dev/claude-code/tiktok-shop-affiliate/logs/launchd_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
